# ansible/roles/dokuwiki/files/supervisor/dokuwiki.conf
[program:dokuwiki]
directory = /opt/dokuwiki
command = docker-compose up --no-color
redirect_stderr = true
autorestart = true
