# ansible/roles/liquid-core/templates/supervisor/liquid-core.conf
[program:liquid-core]
command = /opt/liquid-core/libexec/runserver
