# ansible/roles/liquid-core/templates/supervisor/liquid-core.conf
[program:liquid-core]
user = liquid-apps
command = /opt/liquid-core/libexec/runserver
redirect_stderr = true
autostart = false
