# ansible/roles/network/files/lan/supervisor/lan.conf
[program:lan-dnsmasq]
command = /usr/sbin/dnsmasq -C /var/lib/liquid/lan/dnsmasq.conf
startsecs = 3
redirect_stderr = true
