# ansible/roles/hypothesis/templates/supervisor/hypothesis.conf
[program:hypothesis-elasticsearch]
user = liquid
command = /opt/hypothesis/libexec/elasticsearch
autostart = false
startsecs = 10
