# ansible/roles/hypothesis/templates/supervisor/hypothesis.conf
[program:hypothesis-elasticsearch]
user = liquid
command = /opt/hypothesis/libexec/elasticsearch
autostart = false
startsecs = 10

[program:hypothesis-web]
user = liquid
command = /opt/hypothesis/libexec/run-h gunicorn --name web --paste conf/app.ini -b localhost:46124
autostart = false

[program:hypothesis-websocket]
user = liquid
environment = GUNICORN_STATS_DISABLE="1"
command = /opt/hypothesis/libexec/run-h gunicorn --paste conf/websocket.ini --name websocket --worker-connections 4096 -b localhost:46125
autostart = false

[program:hypothesis-worker]
user = liquid
command = /opt/hypothesis/libexec/run-h hypothesis celery worker
autostart = false

[program:hypothesis-beat]
user = liquid
command = /opt/hypothesis/libexec/run-h hypothesis celery beat
autostart = false
