# ansible/roles/discover/templates/supervisor/avahi-daemon.conf
[program:avahi-daemon]
user = root
command = avahi-daemon -f /etc/avahi/avahi-daemon.conf
redirect_stderr = true
