# ansible/roles/discover/templates/supervisor/discover.conf
[program:liquid-discover]
user = root
command = /opt/discover/libexec/runserver
redirect_stderr = true
