# ansible/roles/discover/templates/supervisor/dnsmasq-dns.conf
[program:dnsmasq-dns]
user = root
command = dnsmasq -k -C /var/lib/liquid/conf/discover/dnsmasq.conf
redirect_stderr = true
autostart = false
startsecs = 2
autorestart = true
stopwaitsecs = 2
