# ansible/roles/matrix/files/supervisor/matrix.conf
[program:matrix]
directory = /opt/matrix
command = docker-compose up --no-color
redirect_stderr = true
