# ansible/roles/dokuwiki/templates/supervisor/dokuwiki.conf
[program:dokuwiki]
directory = /opt/dokuwiki
command = docker-compose up
redirect_stderr = true
