# ansible/roles/network/files/vpn-client/supervisor/vpn-client.conf
[program:vpn-client]
user = root
command = /usr/sbin/openvpn --config /var/lib/liquid/vpn/client.ovpn
redirect_stderr = true
autostart = false
