# ansible/roles/davros/files/supervisor/davros.conf
[program:davros]
directory = /opt/davros
command = docker-compose up --no-color
redirect_stderr = true
