# ansible/roles/network/templates/vpn-server/supervisor/vpn-server.conf
[program:vpn-server]
user = root
command = /usr/sbin/openvpn --config /var/lib/liquid/vpn/server.ovpn
redirect_stderr = true
autostart = false
